# templates/boot/rc.local.erb
#/bin/ksh
#script générique:		
#case $i in		
#	start) ………………….. # Démarrage de mon programm: ex: /usr/bin/su - compte -c "ligne de commande"
#	;;	
#	stop) ……………………. # Arrêt de mon programme	
#	;;	
#	*)	
#	echo 'Usage: /etc/init.d/helloworld {start|stop}'	
#	exit 1	
#	;;	
#esac		